# config/obelisk_decept.conf
start on runlevel [2345]
stop on runlevel [016]

respawn
setuid root
exec /usr/bin/python /opt/obelisk_decept/odlauncher.py
